# files/restic-backup/server-backup-restic-prune-snapshots.timer
[Unit]
Description=Run Restic prune old snapshots once a day around 05:00

[Timer]
OnCalendar=05:10:00
RandomizedDelaySec=20m
Persistent=true
Unit=server-backup-restic-prune-snapshots.service

[Install]
WantedBy=timers.target
