# files/tailscale/tailscale-optimizations.service
[Unit]
Description=tailscale-optimizations
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
ExecStart=/usr/bin/bash -c \
"/usr/sbin/ethtool -K \
$(/usr/sbin/ip -o route get 8.8.8.8 | cut -f 5 -d ' ') \
rx-udp-gro-forwarding on rx-gro-list off"

[Install]
WantedBy=multi-user.target default.target
