# templates/ntfy-server-boot.timer
[Unit]
Description=ntfy-server-boot timer

[Timer]
Unit=ntfy-server-boot.service
OnBootSec=30

[Install]
WantedBy=timers.target
